# thumbor/thumbor-6.0.1/debian/thumbor.service
[Unit]
Description=Thumbor image manipulation service

[Service]
Type=simple
ExecStart=/usr/bin/thumbor
User=thumbor

[Install]
WantedBy=default.target

